# admin-service/admin-service.service
[Unit]
Description=OK TOURS chat-driven admin service
After=network.target

[Service]
Type=simple
User=oktours-admin
Group=oktours-admin
WorkingDirectory=/srv/oktours-repo/admin-service
EnvironmentFile=/etc/oktours-admin.env
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
# /var/www (not just the oktours* subdirs) — publish atomically renames
# directories inside it, which modifies /var/www itself.
ReadWritePaths=/srv/oktours-repo /srv/admin-service /var/www
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
